# playbooks/roles/celery/templates/project.conf
[unix_http_server]
file={{ app_run_dir }}/supervisor.sock  ; (the path to the socket file)
;chmod=0700                             ; socket file mode (default 0700)
;chown=nobody:nogroup                   ; socket file uid:gid owner
;username=user                          ; (default is no username (open server))
;password=123                           ; (default is no password (open server))

[supervisord]
logfile={{ app_log_dir }}/supervisord.log   ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                          ; (num of main logfile rotation backups;default 10)
loglevel=info                               ; (log level;default info; others: debug,warn,trace)
pidfile={{ app_run_dir }}/supervisord.pid   ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                              ; (start in foreground if true;default false)
minfds=1024                                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///{{ app_run_dir }}/supervisor.sock ; use a unix:// URL  for a unix socket

[program:celery_worker]
command={{ app_venv_dir }}/bin/celery -A project worker -E -l info
directory={{ app_root_dir }}
user={{ celery_user }}
autostart=true
autorestart=true
stderr_logfile={{ app_log_dir }}/celery.err.log
stdout_logfile={{ app_log_dir }}/celery.out.log

[program:celery_beat]
command={{ app_venv_dir }}/bin/celery -A project beat -l info
directory={{ app_root_dir }}
autostart=true
autorestart=true
stderr_logfile={{ app_log_dir }}/celery_beat.err.log
stdout_logfile={{ app_log_dir }}/celery_beat.out.log

[program:asgi_daphne]
command={{ app_venv_dir }}/bin/daphne -u {{ app_run_dir }}/daphne.sock project.asgi:application
directory={{ app_root_dir }}
autostart=true
autorestart=true
stderr_logfile={{ app_log_dir }}/daphne.err.log
stdout_logfile={{ app_log_dir }}/daphne.out.log

[program:asgi_worker]
command={{ app_venv_dir }}/bin/python manage.py runworker project
directory={{ app_root_dir }}
process_name=asgi_worker%(process_num)s
numprocs=2
autostart=true
autorestart=true
stderr_logfile={{ app_log_dir }}/runworker.err.log
stdout_logfile={{ app_log_dir }}/runworker.out.log

